# mataroa-notifications.timer
[Unit]
Description=Run mataroa-notifications every day at 10AM UTC

[Timer]
OnCalendar=*-*-* 10:00:00

[Install]
WantedBy=timers.target
