# ansible/mataroa-exports.service.j2
[Unit]
Description=Send mataroa exports

[Service]
Type=oneshot
User=deploy
WorkingDirectory=/var/www/mataroa
EnvironmentFile=/etc/systemd/system/mataroa.env
ExecStart=/home/deploy/.local/bin/uv run manage.py mailexports

[Install]
WantedBy=multi-user.target
