# ansible/mataroa-backup.service.j2
[Unit]
Description=Backup mataroa database

[Service]
Type=oneshot
User=deploy
WorkingDirectory=/home/deploy
Environment="PGPASSWORD={{ postgres_password }}"
ExecStart=/home/deploy/backup-database.sh

[Install]
WantedBy=multi-user.target
